# chef-repo/wrapper_cookbooks/tubmanproject_jupyterhub/templates/default/supervisor_program.conf.erb
[program:<%= @program_name%>]
<% @program_config.each_pair do |key, value| %>
<%= "#{key}=#{value}"%>
<% end %>